# .config/systemd/user/bluetooth-autoconnect.service
[Unit]
Description=Bluetooth autoconnect service
Before=bluetooth.service

[Service]
Type=simple
ExecStart=/home/serge/bin/bluetooth-autoconnect -d

[Install]
WantedBy=bluetooth.service
