# .config/systemd/user/dunst.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dunst -verbosity info -print
Restart=always
RestartSec=1s
RuntimeMaxSec=3600s
