# .config/systemd/user/dunst.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dunst -verbosity info
Restart=
Restart=always
RestartSec=
RestartSec=1s

