# .config/systemd/user/dunst.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dunst -verbosity info
